# mongodb/mongo.crontab

30 1 * * * bash /opt/aiddata/backup_cron_log.sh master
15 1 * * * bash /opt/aiddata/backup_mongo_log.sh master 2>&1 | tee 1>>/opt/aiddata/cron_log.log
0 */6 * * * bash /opt/aiddata/backup_mongo_db.sh master 2>&1 | tee 1>>/opt/aiddata/cron_log.log
