# mongodb/mongo.crontab

1 0 * * * bash /opt/aiddata/backup_cron_log.sh master
1 5 * * * bash /opt/aiddata/backup_mongo_log.sh master 2>&1 | tee 1>>/opt/aiddata/cron_log.log
1 10 * * * bash /opt/aiddata/backup_mongo_db.sh master 2>&1 | tee 1>>/opt/aiddata/cron_log.log
